# desktop_agent/templates/npamx-agent.service
[Unit]
Description=NPAMX Desktop Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/npamx-agent run
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
